# etc/systemd/system/battery-icon.service
[Unit]
Description=Refresh the battery percentage desktop icon (Acoustic Camera UPS)

[Service]
Type=oneshot
User=jdn
WorkingDirectory=/home/jdn/Code/Acoustic-Camera
ExecStart=/home/jdn/Code/Acoustic-Camera/venv/bin/python /home/jdn/Code/Acoustic-Camera/bin/battery_icon.py
